# software/rootfs/etc/init/robot.conf
description 'node upstart script for instances'
author 'Tom Navarro'

start on (local-filesystems and net-device-up)
stop on runlevel [06]

script
	nodemon /home/ubuntu/widgets/app.js
end script
